# etc/systemd/xdg-retroarch-favourites.path
[Unit]
Description=(Trigger for) Sync Retroarch favourites directly into the desktop menu system

[Path]
PathChanged=%h/.config/retroarch/content_favorites.lpl

[Install]
WantedBy=paths.target
